# scripts/install/pkgbuilds/nyxtra-network/nyxtra-network.install
post_install() {
  echo ">>> Enabling NetworkManager..."
  systemctl enable NetworkManager.service || true
}
