# scripts/install/pkgbuilds/nyxtra-audio/nyxtra-audio.install
post_install() {
  echo ">>> Enabling Audio (for all users)..."
  systemctl --global enable --now pipewire.service wireplumber.service pipewire-pulse.service
  systemctl --global enable --now bluetooth.service

}
