# scripts/install/pkgbuilds/nyxtra-audio/nyxtra-audio.install
post_install() {
  echo ">>> Enabling Audio (for all users)..."
  systemctl --global enable pipewire.service wireplumber.service pipewire-pulse.service
}
