# profiles/pkgbuilds/nyxtra-login/nyxtra-login.install
post_install() {
  echo ">>> Enabling Ly greeter..."
  systemctl enable ly@tty1.service || true
}
